# ops/football-predict-live.timer
[Unit]
Description=EPL Prediction Lab — poll live scores every 3s (API-Football Pro)

[Timer]
OnBootSec=3s
OnUnitActiveSec=3s
Unit=football-predict-live.service
AccuracySec=500ms

[Install]
WantedBy=timers.target
