# ops/football-predict-live.timer
[Unit]
Description=EPL Prediction Lab — poll live scores every 1s (API-Football Ultra)

[Timer]
OnBootSec=1s
OnUnitActiveSec=1s
Unit=football-predict-live.service
AccuracySec=100ms

[Install]
WantedBy=timers.target
