# ops/football-predict-live.service
[Unit]
Description=EPL Prediction Lab — live score polling daemon (API-Football)
After=docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/football-predict
ExecStart=/usr/bin/docker compose exec -T api python scripts/ingest_live_scores.py --daemon --interval=1.0
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
