# ops/football-predict-rho.timer
[Unit]
Description=football-predict — 1st of month 05:30 UTC ρ recalibration

[Timer]
OnCalendar=*-*-01 05:30:00
Unit=football-predict-rho.service
Persistent=true

[Install]
WantedBy=timers.target
