# ops/football-predict-live.timer
[Unit]
Description=EPL Prediction Lab — poll live scores every 2s (API-Football Ultra)

[Timer]
OnBootSec=2s
OnUnitActiveSec=2s
Unit=football-predict-live.service
AccuracySec=500ms

[Install]
WantedBy=timers.target
