# ops/football-predict-live.service
[Unit]
Description=EPL Prediction Lab — live score polling (API-Football)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/football-predict
ExecStart=/bin/bash /opt/football-predict/ops/live.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=2m
