# ops/football-predict-live.timer
[Unit]
Description=EPL Prediction Lab — poll live scores every 5s

[Timer]
OnBootSec=5s
OnUnitActiveSec=5s
Unit=football-predict-live.service
AccuracySec=1s

[Install]
WantedBy=timers.target
