# ops/football-predict-news.service
[Unit]
Description=EPL Prediction Lab — RSS news aggregator
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/football-predict
ExecStart=/bin/bash /opt/football-predict/ops/news.sh
User=root
