# ops/football-predict-weekly.timer
[Unit]
Description=Run EPL Prediction Lab weekly refresh every Monday 02:00 UTC

[Timer]
OnCalendar=Mon *-*-* 02:00:00 UTC
Persistent=true
RandomizedDelaySec=10min
Unit=football-predict-weekly.service

[Install]
WantedBy=timers.target
