# ops/football-predict-live.timer
[Unit]
Description=EPL Prediction Lab — poll live scores every 10s

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s
Unit=football-predict-live.service
AccuracySec=2s

[Install]
WantedBy=timers.target
