# ops/football-predict-lineups.timer
[Unit]
Description=football-predict — lineups every 5 minutes (API-Football)

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=football-predict-lineups.service
AccuracySec=30s

[Install]
WantedBy=timers.target
