# ops/football-predict-hourly.timer
[Unit]
Description=football-predict — every 2h freshness refresh

[Timer]
OnCalendar=*-*-* 00,02,04,06,08,10,12,14,16,18,20,22:00:00
Unit=football-predict-hourly.service
Persistent=true
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
